# auth_gate.service
[Unit]
Description=REO Dashboard Authentication Gateway
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/reo-config
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/bin/python3 /var/www/reo-config/auth_gate.py
Restart=always
RestartSec=10
StandardOutput=append:/var/www/iproot/reo/logs/auth_gateway.log
StandardError=append:/var/www/iproot/reo/logs/auth_gateway.log

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/iproot/reo/logs
ReadOnlyPaths=/var/www/reo-config

[Install]
WantedBy=multi-user.target

